# deploy/systemd/hammertrade-paper.example.service
[Unit]
Description=HammerTrade Paper Trader — SiM6 SELL (no real orders)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/hammertrade
EnvironmentFile=/opt/hammertrade/.env
ExecStart=/opt/hammertrade/.venv/bin/python scripts/run_paper_trader.py \
    --ticker SiM6 \
    --class-code SPBFUT \
    --timeframe 1m \
    --profile balanced \
    --direction-filter SELL \
    --env prod \
    --market-hours-config configs/market_hours/moex_futures.yaml \
    --api-timeout-sec 10 \
    --status-file runtime/paper_status_SiM6_SELL.json
Restart=always
RestartSec=10
User=andrey

[Install]
WantedBy=multi-user.target
